# relay/gisgmp/gisgmp-relay.service
[Unit]
Description=Релей ГИС ГМП -> ЖКХ (демон, управление из ЖКХ по pull)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/opt/gisgmp-relay/relay.env
# Маршрут к корп-сети (реестр виден только через Cisco 10.23.0.1).
# '-' и 'replace' — идемпотентно, не валит старт, если уже есть.
ExecStartPre=-/usr/sbin/ip route replace 10.0.0.0/8 via 10.23.0.1 dev ens19
ExecStart=/usr/bin/python3 /opt/gisgmp-relay/relay.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
