# relay/gisgmp/gisgmp-relay.service
[Unit]
Description=Релей ГИС ГМП -> ЖКХ (подгрузка долгов жильцов)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/opt/gisgmp-relay/relay.env
# Маршрут к корп-сети (реестр виден только через Cisco 10.23.0.1).
# '-' и 'replace' — идемпотентно и не валит запуск, если уже есть.
ExecStartPre=-/usr/sbin/ip route replace 10.0.0.0/8 via 10.23.0.1 dev ens19
ExecStart=/usr/bin/python3 /opt/gisgmp-relay/relay.py
# Логи — в journald (journalctl -u gisgmp-relay).
